# thecave/etc/rc.local
#!/bin/bash


if [ -f /aafirstboot ]; then /aafirstboot start ; fi

su - intwocave -c "cp -r /home/intwocave/.git/dgk /home/intwocave/ramdisk"

exit 0
